# yocto/meta-custom/recipes-apps/myspotifydisplay/files/spotifydisplay.service
[Unit]
Description=Spotify Display Application
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/spotifydisplay
Restart=no

[Install]
WantedBy=multi-user.target
